# crontab
*/3 * * * * python src/main.py >> ./logfile